# src/ClusterBootstrap/deploy/storage/auto_share/glusterfs.mount
[Unit]
Description=Mount GlusterFS volume automatically  
After=docker.service


[Mount]
What=:
Where=
Type=glusterfs
Options=
